# root/etc/s6-overlay/s6-rc.d/pyload-cleanup-pid/run
#!/command/with-contenv sh

if [ -e /config/pyload.pid ]; then
  rm /config/pyload.pid
fi

