# root/etc/s6-overlay/s6-rc.d/pyload-cleanup-pid/up
/etc/s6-overlay/s6-rc.d/pyload-cleanup-pid/run
